# deploy/amazon-us-scraper.service.example
# 复制到服务器后：
#   sudo cp deploy/amazon-us-scraper.service.example /etc/systemd/system/amazon-us-scraper.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now amazon-us-scraper
#
# 约定：监听 0.0.0.0:8989（云安全组需放行 TCP 8989，若仅本机访问可改为 127.0.0.1:8989）

[Unit]
Description=Amazon US Scraper (FastAPI + Uvicorn)
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/amazon-us-scraper
Environment="PATH=/opt/amazon-us-scraper/.venv/bin"
# 环境变量由应用内 python-dotenv 加载项目根 .env（避免 systemd 解析 .env 与注释/编码不一致）
ExecStart=/opt/amazon-us-scraper/.venv/bin/uvicorn webapp.main:app --host 0.0.0.0 --port 8989
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
